# deploy/systemd/radcast.service
[Unit]
Description=RADcast production service
After=network.target

[Service]
Type=simple
User=fldadmin
Group=fldadmin
WorkingDirectory=/home/fldadmin/radcast
EnvironmentFile=/home/fldadmin/radcast/.env
ExecStart=/home/fldadmin/radcast/.venv/bin/gunicorn -k uvicorn.workers.UvicornWorker -w 1 -b 127.0.0.1:8013 radcast.api:app
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
